# root/etc/s6-overlay/s6-rc.d/init-prowlarr-config/run
#!/usr/bin/with-contenv bash

# permissions
lsiown -R abc:abc \
    /config
