# root/etc/s6-overlay/s6-rc.d/init-prowlarr-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

mkdir -p /run/prowlarr-temp

# permissions
lsiown -R abc:abc \
    /config \
    /run/prowlarr-temp
